# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=Xvfb :0 -screen 0 1024x768x16
autorestart=true
priority=100

[program:fluxbox]
command=fluxbox
environment=DISPLAY=":0"
autorestart=true
priority=200

[program:x11vnc]
command=x11vnc -display :0 -forever -nopw -listen localhost -rfbport 5900
autorestart=true
priority=300

[program:novnc]
command=websockify --web /usr/share/novnc %(ENV_PORT)s localhost:5900
autorestart=true
priority=400

[program:chromium]
# --no-sandbox is MANDATORY inside Docker/Heroku
# --disable-dev-shm-usage prevents crashes in low-memory environments
command=chromium-browser --display=:0 --no-sandbox --disable-dev-shm-usage --disable-gpu --start-maximized
autorestart=true
priority=500
