# supervisord.conf
[supervisord]
nodaemon=true
# Removed user=root to allow Heroku to run the process as a non-privileged user
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=Xvfb :0 -screen 0 1280x720x24
autorestart=true
priority=10

[program:fluxbox]
command=fluxbox
autorestart=true
priority=20
environment=DISPLAY=":0"

[program:x11vnc]
command=x11vnc -display :0 -forever -shared -rfbport 5900 -nopw
autorestart=true
priority=30

[program:novnc]
# Ubuntu's 'novnc' package uses websockify to serve the files in /usr/share/novnc
command=websockify --web=/usr/share/novnc/ %(ENV_PORT)s localhost:5900
autorestart=true
priority=40

[program:firefox]
# Ensure the profile directory is created by the current user and cleaned before launch
command=/bin/bash -c "mkdir -p /tmp/firefox-profile && rm -f /tmp/firefox-profile/.parentlock && firefox --display=:0 --profile /tmp/firefox-profile --no-remote"
autorestart=true
priority=50
environment=DISPLAY=":0"
