# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:xvfb]
# Create the virtual display :0
command=Xvfb :0 -screen 0 1024x768x16
autorestart=true
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fluxbox]
# The window manager so you can actually move windows around
command=fluxbox
environment=DISPLAY=":0"
autorestart=true
priority=200
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
# The VNC server that captures the Xvfb screen
command=x11vnc -display :0 -forever -nopw -listen localhost -rfbport 5900 -o /tmp/x11vnc.log
autorestart=true
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
# This bridges the VNC traffic to the Heroku $PORT via WebSockets
command=websockify --web /usr/share/novnc %(ENV_PORT)s localhost:5900
autorestart=true
priority=400
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:chromium]
# Start Chromium with Docker/Heroku compatibility flags
command=chromium-browser --display=:0 --no-sandbox --disable-dev-shm-usage --disable-gpu --remote-debugging-port=9222
autorestart=true
priority=500
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
