# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:setup-extension]
# This program creates the extension folder and files
command=bash -c "mkdir -p /tmp/my-extension && echo '{\"manifest_version\": 3, \"name\": \"My Custom Extension\", \"version\": \"1.0\", \"description\": \"Pre-loaded unpacked extension\", \"permissions\": [\"storage\"], \"action\": { \"default_title\": \"Hello!\" }}' > /tmp/my-extension/manifest.json && echo 'console.log(\"Extension Loaded!\");' > /tmp/my-extension/background.js"
autorestart=false
startsecs=0
priority=50

[program:xvfb]
command=Xvfb :0 -screen 0 1024x768x16
autorestart=true
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fluxbox]
command=fluxbox
environment=DISPLAY=":0"
autorestart=true
priority=200
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
command=x11vnc -display :0 -forever -nopw -listen localhost -rfbport 5900 -o /tmp/x11vnc.log
autorestart=true
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
command=websockify --web /usr/share/novnc %(ENV_PORT)s localhost:5900
autorestart=true
priority=400
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:chrome]
# Added --test-type to hide the "unsupported flag" warning
# Added --disable-blink-features=AutomationControlled to bypass Cloudflare/Bot detection
# Added a standard User-Agent and window-size to look more like a real user
command=google-chrome-stable \
  --no-sandbox \
  --test-type \
  --disable-blink-features=AutomationControlled \
  --user-agent="Mozilla/5.0 (Windows NT 10.0; Win64; x64) AppleWebKit/537.36 (KHTML, like Gecko) Chrome/131.0.0.0 Safari/537.36" \
  --window-size=1920,1080 \
  --start-maximized \
  --disable-gpu \
  --disable-software-rasterizer \
  --disable-gpu-compositing \
  --use-gl=swiftshader \
  --disable-dev-shm-usage \
  --user-data-dir=/tmp/chrome-user-data \
  --load-extension=/tmp/my-extension \
  --enable-extensions \
  --no-first-run \
  --no-default-browser-check \
  --disable-background-networking \
  --disable-background-timer-throttling \
  --disable-breakpad \
  --disable-client-side-phishing-detection \
  --disable-component-update \
  --disable-domain-reliability \
  --disable-features=TranslateUI,GCM \
  --disable-sync \
  --disable-translate \
  --disable-extensions-file-access-check \
  --disable-dbus \
  --disable-notifications \
  --disable-logging \
  --log-level=3 \
  --display=:0
autorestart=true
priority=500
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
