# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:xvfb]
# Create the virtual display :0
command=Xvfb :0 -screen 0 1024x768x16
autorestart=true
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fluxbox]
# The window manager so you can actually move windows around
command=fluxbox
environment=DISPLAY=":0"
autorestart=true
priority=200
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
# The VNC server that captures the Xvfb screen
# Added -o /tmp/x11vnc.log to ensure the log is writable
command=x11vnc -display :0 -forever -nopw -listen localhost -rfbport 5900 -o /tmp/x11vnc.log
autorestart=true
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
# This bridges the VNC traffic to the Heroku $PORT via WebSockets
# %(ENV_PORT)s is how Supervisor reads the $PORT variable from Heroku
command=websockify --web /usr/share/novnc %(ENV_PORT)s localhost:5900
autorestart=true
priority=400
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:chrome]
# (Optional) Auto-starts Google Chrome so the screen isn't just a blank gray void
# Flags:
#  --no-sandbox kept for container root execution
#  --user-data-dir provides a writable profile so extensions can be installed and accessed
#  --enable-extensions ensures extension system is enabled
#  --disable-dev-shm-usage avoids /dev/shm issues in containers
#  --no-first-run and --no-default-browser-check reduce first-run prompts
command=google-chrome-stable --no-sandbox --disable-dev-shm-usage --user-data-dir=/tmp/chrome-user-data --enable-extensions --no-first-run --no-default-browser-check --display=:0
autorestart=true
priority=500
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
